# crontab
0 * * * * cd ~/Repos/ISP-Speed && ./generate_report.sh