# speed_report/crontab
0 * * * * cd ~/Downloads/speed_report && ./generate_report.sh
