# speed_report/crontab
0 * * * * cd ~/Repos/Every-Day-Technology/speed_report && ./generate_report.sh
